# templates/restic.timer.j2
[Unit]
Description=Run restic every day

[Timer]
OnCalendar=3:0:00

[Install]
WantedBy=timers.target